# pcsd/pcsd.service.in
[Unit]
Description=PCS GUI and remote configuration interface
Documentation=man:pcsd(8)
Documentation=man:pcs(8)
Requires=network-online.target
After=network-online.target
Requires=pcsd-ruby.service
After=pcsd-ruby.service

[Service]
EnvironmentFile=@CONF_DIR@/pcsd
ExecStart=@SBINDIR@/pcsd
StateDirectory=pcsd
StateDirectoryMode=0700
LogsDirectory=pcsd
LogsDirectoryMode=0700
Type=notify
KillMode=mixed

[Install]
WantedBy=multi-user.target
