# pcsd/pcsd-ruby.service.in
[Unit]
Description=PCS GUI and remote configuration interface (Ruby)
Documentation=man:pcsd(8)
Documentation=man:pcs(8)
Requires=network-online.target
After=network-online.target
# Stop the service automatically if nothing that depends on it is running
StopWhenUnneeded=true
# When stopping or restarting pcsd, stop or restart pcsd-ruby as well
PartOf=pcsd.service

[Service]
EnvironmentFile=@CONF_DIR@/pcsd
@SYSTEMD_GEM_HOME@
# This file holds the selinux context
ExecStart=@LIB_DIR@/pcsd/pcsd
StateDirectory=pcsd
StateDirectoryMode=0700
LogsDirectory=pcsd
LogsDirectoryMode=0700
Type=notify

[Install]
WantedBy=multi-user.target
