# pcsd/pcsd.service
[Unit]
Description=PCS GUI and remote configuration interface

[Service]
EnvironmentFile=-/etc/sysconfig/pcsd
WorkingDirectory=/var/lib/pcsd
ExecStart=/usr/lib/pcsd/pcsd start
Type=notify
NotifyAccess=all

[Install]
WantedBy=multi-user.target
