# pcsd/pcsd.service
[Unit]
Description=PCS GUI and remote configuration interface
Documentation=man:pcsd(8)
Documentation=man:pcs(8)

[Service]
EnvironmentFile=/etc/sysconfig/pcsd
ExecStart=/usr/sbin/pcsd
Type=notify

[Install]
WantedBy=multi-user.target
